# jeebs.service
#
# This is a TEMPLATE file for reference.
# It is STRONGLY recommended to use the `install.sh` script to generate and install
# the service file, as it will automatically fill in the correct user and paths for your VPS.
#
[Unit]
Description=JeebsAI Server
After=network.target

[Service]
Type=simple
# Run as a standard user since Nginx handles port 80/443.
# Replace 'your_username' with your actual Linux username.
User=your_username
WorkingDirectory=/path/to/JeebsAI-1
ExecStart=/path/to/JeebsAI-1/target/release/jeebs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target